# scripts/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:php-fpm]
directory=/tmp
home=/usr/share/nginx/html
command=php-fpm7 -F
priority=30
autostart=true
autorestart=false
stopsignal=QUIT
stdout_logfile=/tmp/phpfpm.stdout
stderr_logfile=/tmp/phpfpm.stderr

[program:nginx]
directory=/etc/nginx
home=/etc/nginx
command=/usr/sbin/nginx -g "daemon off;" -c /scripts/nginx.conf
priority=40
autostart=true
autorestart=false
stdout_logfile=/tmp/nginx.stdout
stderr_logfile=/tmp/nginx.stderr
